# vg-sentry.service
[Unit]
Description=Volantic Genesis background health check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/vg health

[Install]
WantedBy=default.target
