# legacy_python/genesis-sentry.service
[Unit]
Description=Genesis Sentry background health check

[Service]
Type=oneshot
ExecStart=/usr/local/bin/genesis monitor