# legacy_python/genesis-greet.service
[Unit]
Description=Genesis startup greeting

[Service]
Type=oneshot
ExecStart=/usr/local/bin/genesis greet

[Install]
WantedBy=default.target
